# deploy/systemd/vpn-web-admin.service
[Unit]
Description=MOROZ VPN Web Admin (FastAPI)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/vpn-bot
Environment=PYTHONPATH=/opt/vpn-bot
ExecStart=/usr/bin/python3 -m uvicorn web_admin.backend.main:app --host 0.0.0.0 --port 8889
Restart=always
RestartSec=5
User=www-data

[Install]
WantedBy=multi-user.target

