# deploy/systemd/vpn-web-admin.service
[Unit]
Description=VPN Bot Web Admin Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/vpn-bot
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=/opt/vpn-bot/.env
ExecStart=/usr/bin/python3 /opt/vpn-bot/web_admin/backend/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vpn-web-admin

[Install]
WantedBy=multi-user.target

