# yoomoney_backend/yoomoney-backend.service
[Unit]
Description=YooMoney Backend Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/yoomoney_backend
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /root/yoomoney_backend/app.py
Restart=always
RestartSec=10
StandardOutput=append:/root/yoomoney_backend/app.log
StandardError=append:/root/yoomoney_backend/app.log

[Install]
WantedBy=multi-user.target

